# javaserver.sv.conf
[program:javaserver]
command=java -jar /app/server.jar
user=root
environment=DISPLAY=':1.5'
stdout_logfile=/tmp/javaserver.log
stderr_logfile=/tmp/javaserver.err.log
